# conf/systemd/stop.service
[Unit]
Description=stop: terminates systemd

[Service]
Type=oneshot
#ExecStart=/bin/kill -2 1
ExecStart=/bin/true
